# packages/net-im/prosody/files/systemd/prosody.service
[Unit]
Description=Prosody XMPP server
After=syslog.target

[Service]
Type=simple
User=prosody
Group=prosody
ExecStart=/usr/bin/prosody

[Install]
WantedBy=multi-user.target
